# ha-dispatch/rootfs/etc/services.d/ha-dispatch/finish
#!/command/execlineb -S0
# On service exit, bring the whole container down (HA will restart it).
s6-svscanctl -t /var/run/s6/services
